# aria2.service.template
[Unit]
Description=aria2
After=

[Service]
Type=forking
ExecStart= template/run_aria2.sh
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always
PrivateTmp=true

[Install]
WantedBy=multi-user.target

